# assets/nu.wigge.smedja.smdjad.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  TEMPLATE — NOT FOR DIRECT USE.

  The PLACEHOLDER_INSTALL_DIR / PLACEHOLDER_HOME tokens below are never
  substituted at runtime, so copying this file verbatim into
  ~/Library/LaunchAgents will produce a broken launchd job.

  This file is reference only and is deliberately NOT shipped in the release
  tarball. The real per-user plist is generated, fully resolved, by
  scripts/install.sh (see its launchd setup section) on macOS installs.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>nu.wigge.smedja.smdjad</string>

  <key>ProgramArguments</key>
  <array>
    <string>PLACEHOLDER_INSTALL_DIR/smdjad</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>PLACEHOLDER_INSTALL_DIR:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>XDG_RUNTIME_DIR</key>
    <string>/tmp</string>
  </dict>

  <key>StandardOutPath</key>
  <string>PLACEHOLDER_HOME/Library/Logs/smedja/smdjad.log</string>

  <key>StandardErrorPath</key>
  <string>PLACEHOLDER_HOME/Library/Logs/smedja/smdjad-error.log</string>

  <key>WorkingDirectory</key>
  <string>PLACEHOLDER_HOME</string>
</dict>
</plist>
