# assets/smdjad.service
[Unit]
Description=smedja AI orchestration daemon
After=network.target graphical-session.target

[Service]
# Type=notify: smdjad sends sd_notify(READY=1) once the control socket is bound
# and the database is open, so systemd only considers the unit started when it
# is actually ready to serve.
Type=notify
ExecStartPre=-/bin/mkdir -p %h/.config/smedja %h/.local/share/smedja
# Provider API keys pasted via the TUI `/login` masked entry are written here
# (0600). The leading `-` makes the file optional, so the daemon still starts
# when no keys have been saved.
EnvironmentFile=-%h/.config/smedja/secrets.env
ExecStart=%h/.local/bin/smdjad
Restart=on-failure
RestartSec=5s
TimeoutStopSec=35
StandardOutput=journal
StandardError=journal
KillMode=mixed
Environment=XDG_RUNTIME_DIR=/run/user/%U
# Runner CLIs (codex, opencode, npm-global tools) commonly install under the
# user's home bins, which a systemd user service does NOT inherit. Put them on
# PATH so the provider pool can detect them at startup.
Environment=PATH=%h/.local/bin:%h/.npm-global/bin:%h/bin:/usr/local/bin:/usr/bin:/bin

# ── Sandboxing ────────────────────────────────────────────────────────────
# smdjad's whole job is to run coding agents (claude-cli, codex, opencode) and
# their tools. Those children must read/run/EDIT in the user's workspace and
# must initialise their OWN command sandboxes (claude-cli sandboxes Bash via a
# user-namespace/bubblewrap setup). systemd unit hardening is inherited by every
# child, so an aggressive profile (ProtectSystem=strict + ProtectHome=read-only
# + NoNewPrivileges=yes) breaks them with "sandbox setup fails on a read-only
# filesystem" — the agent can't run a single command.
#
# Command isolation is therefore delegated to the layer designed for it:
# smedja's per-tool sandbox (landlock/seatbelt/docker) and the provider's own
# sandbox. The unit keeps only protection that does NOT strangle children:
#   * ProtectSystem=full — /usr, /boot, /efi, /etc stay read-only; the user's
#     workspace under $HOME stays writable so agents can read/edit code.
#   * NoNewPrivileges=no — required so nested agent sandboxes can initialise.
#   * PrivateTmp=yes — private, writable /tmp for the daemon process tree.
PrivateTmp=yes
ProtectSystem=full
NoNewPrivileges=no
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=default.target
