# assets/smdjad.service
[Unit]
Description=smedja AI orchestration daemon
After=network.target graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/smdjad
Restart=on-failure
RestartSec=5s
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
